# docker/configs/supervisor-app.conf
[program:app-uwsgi]
command = /usr/local/bin/uwsgi --ini /root/DVA/docker/configs/uwsgi.ini

[program:nginx-app]
command = /usr/sbin/nginx